# wikum/supervisord.conf
[unix_http_server]
file=/wikum/wikum/wikum/supervisor.sock   ; path to your socket file

[supervisord]
logfile=/var/log/supervisord/supervisord.log ; supervisord log file
logfile_maxbytes=50MB       ; maximum size of logfile before rotation
logfile_backups=10          ; number of backed up logfiles
loglevel=info               ; info, debug, warn, trace
pidfile=/tmp/supervisord.pid ; pidfile location
nodaemon=false              ; run supervisord as a daemon
minfds=1024                 ; number of startup file descriptors
minprocs=200                ; number of process descriptors
childlogdir=/var/log/supervisord/            ; where child log files will live

[program:celeryd]
user=ubuntu
command=/wikum-env3/bin/celery -A wikum.celery worker -l info -c 2
numprocs=2
process_name = %(program_name)s_%(process_num)02d
stdout_logfile=/var/log/celeryd.log
stderr_logfile=/var/log/celeryd.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
priority=998

[program:asgi_daphne]
command=/wikum-env3/bin/daphne -b 0.0.0.0 -p 8001 wikum.asgi:application
directory=/wikum/wikum/wikum
user=ubuntu
numprocs=1
stdout_logfile=/var/log/stdout_daphne.log
stderr_logfile=/var/log/stderr_daphne.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
killasgroup=true
priority=998

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///wikum/wikum/wikum/supervisor.sock ; use unix:// schem for a unix sockets.

